# lib/systemd/system/offlineimap.timer
[Unit]
Description=Offlineimap Main Timer

[Timer]
OnBootSec=1m
OnUnitInactiveSec=5m

[Install]
WantedBy=default.target
